# deploy/supervisor/mayoreo-cloud.conf
; Supervisor: worker de cola Laravel para mayoreo-cloud
; Copiar a: /etc/supervisor/conf.d/mayoreo-cloud.conf
; Ajustar 'user' y rutas si no usas root /home/mayoreo/htdocs/mayoreo-cloud

[program:mayoreo-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/mayoreo/htdocs/mayoreo-cloud/artisan queue:work database --queue=high,default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/home/mayoreo/htdocs/mayoreo-cloud/storage/logs/worker.log
stopwaitsecs=3600
