# deploy/supervisor-laravel-queue.conf
; Worker de cola Laravel para mayoreo-cloud (envío de ofertas a Telegram).
; Si el proyecto está en otra ruta, edita command, stdout_logfile y directory.
; Copia a /etc/supervisor/conf.d/ y ejecuta: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start laravel-queue:*
; Ver docs/COLA-WORKER-AUTOMATIZAR.md para más detalles.

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /home/mayoreo/htdocs/mayoreo-cloud/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/home/mayoreo/htdocs/mayoreo-cloud/storage/logs/worker.log
stopwaitsecs=60
directory=/home/mayoreo/htdocs/mayoreo-cloud
